# mac-sync-bridge/launchd/com.iosgtd.syncbridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.iosgtd.syncbridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>swift</string>
    <string>run</string>
    <string>BridgeApp</string>
    <string>--config</string>
    <string>__CONFIG_PATH__</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__WORKDIR__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/bridge.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/bridge.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>
</dict>
</plist>
